# root/etc/services.d/nfs/run
#!/usr/bin/with-contenv bash

exec inotifywait -m `/nfs-mount-targets`